# config/systemd/user/aria2.service
[Unit]
Description=Aria2

[Service]
ExecStart=/usr/bin/aria2c --log=-
Restart=always

[Install]
WantedBy=default.target